# systemd/doc-store-filewatcher.service
[Unit]
Description=doc-store file watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=doc-store
Group=doc-store
WorkingDirectory=/opt/doc-store
Environment=DOC_STORE_FILEWATCHER_CONFIG=/etc/doc-store/filewatcher.yaml
ExecStart=/opt/doc-store/.venv/bin/python -m doc_store_filewatcher.service
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
